# macos/launchd/com.rakshak.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rakshak.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/rakshak/bin/RakshakDaemon</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RAKSHAK_EXTERNAL_COREDNS</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Library/Logs/Rakshak/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/Rakshak/daemon.error.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/rakshak</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
